# etc/service/omoide-worker.service
[Unit]
Description=Omoide worker daemon
After=network.target

[Service]
User=omoide-user
WorkingDirectory=/home/omoide-user/omoide
ExecStart=/home/omoide-user/omoide/run_worker.sh
EnvironmentFile=/etc/omoide/env
CPUAccounting=true
MemoryAccounting=true

[Install]
WantedBy=multi-user.target
